# resque/templates/default/process.conf.erb
[program:<%= @name %>]
command=<%= @command %>
process_name=<%= @name %>-%(process_num)s
numprocs=<%= @number %>
user=<%= @user %>
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/<%= @name %>-%(process_num)s.stdout.log
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/<%= @name %>-%(process_num)s.stderr.log
stderr_logfile_backups=3
environment=<%= @env.select{ |key, value| value }.map{ |key, value| "#{key}=\"#{value}\"" }.join(',') %>
